# script/keymapSetting/mapLoader.service
[Unit]
Description=mapLoader

[Service]
ExecStart=/usr/bin/loadkeys colemak

[Install]
WantedBy=multi-user.target
